# lib/base/lib/assets/sync-checker/sync-checker.service
[Unit]
Description="Sync checker for blockchain node"

[Service]
ExecStart=/opt/syncchecker.sh
