# lib/solana/lib/assets/instance/network/net-sync-checker.service
[Unit]
Description="Net sync checker for blockchain node"

[Service]
ExecStart=/opt/instance/network/net-syncchecker.sh
