# lib/solana/lib/assets/systemd/sync-checker.service
[Unit]
Description="Sync checker for the tezos node"

[Service]
ExecStart=/opt/syncchecker.sh