# lib/solana/lib/assets/systemd/sync-checker.timer
[Unit]
Description="Run Sync checker service every 5 min"

[Timer]
OnCalendar=*:*:0/5
Unit=sync-checker.service

[Install]
WantedBy=multi-user.target